# utils/utils.crontab
# Do NOT:
#   crontab -e
# Do:
#   crontab utils.crontab

# healthchecks.io checks are supposed to be in 10-minute interval
*/10 * * * * cd $HOME/utils && . ./.env && ./healthchecks.sh
